# deploy/gamehub.service
[Unit]
Description=GAMEHUB — family LAN game library (port 8096)
After=network-online.target

[Service]
WorkingDirectory=%h/projects/gamehub
ExecStart=%h/projects/gamehub/.venv/bin/python %h/projects/gamehub/server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
